# static/etc/services.d/resilio/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid bytesized /app/rslsync --nodaemon --config /config/sync.conf
